# prestodb/cdh5.12-hive/files/etc/supervisord.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
startsecs=2
stopwaitsecs=10
startretries=30
user=root
redirect_stderr=true
stdout_logfile=/var/log/sshd
autostart=true
