# prestodev/cdh5.12-hive/files/etc/supervisord.d/hive-metastore.conf
[program:hive-metastore]
command=hive --service metastore
startsecs=2
stopwaitsecs=10
user=hive
redirect_stderr=true
stdout_logfile=/var/log/hive/hive-metastore.log
autostart=true
