# teradatalabs/cdh5-hive/files/supervisord.d/sshd.conf
[program:sshd]
command=/etc/init.d/sshd start
priority=1
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/sshd.start
autostart=false
