# prestosql/hdp2.5-hive/files/supervisord.d/zookeeper.conf
[program:zookeeper]
command=/usr/bin/zookeeper-server start-foreground
startsecs=2
stopwaitsecs=10
user=zookeeper
redirect_stderr=true
stdout_logfile=/var/log/zookeeper/zookeeper.log
autostart=false
