# testing/accumulo/files/etc/supervisord.d/accumulo-manager.conf
[program:accumulo-master]
command=/usr/local/lib/accumulo/bin/accumulo manager
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/accumulo/manager.log
autostart=true
