# archived/accumulo/files/etc/supervisord.d/accumulo-monitor.conf
[program:accumulo-monitor]
command=/usr/local/lib/accumulo/bin/accumulo monitor
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/accumulo/monitor.log
autostart=true
