# prestosql/hdp2.5-hive-kerberized/files/supervisord.d/bootstrap-kerberos.conf
[program:bootstrap-kerberos]
command=/root/bootstrap-kerberos.sh
priority=0
redirect_stderr=true
stdout_logfile=/var/log/bootstrap-kerberos.log
autostart=true
autorestart=false
