# prestodev/hive3.1-hive/files/etc/supervisord.d/hdfs-namenode.conf
[program:hdfs-namenode]
command=hdfs namenode
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/hadoop-hdfs/hadoop-hdfs-namenode.log
autostart=true
