# prestodev/hive3.1-hive/files/etc/supervisord.d/hive-server2.conf
[program:hive-server2]
command=hive --service hiveserver2
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/hive/hive-server2.log
autostart=true
