# commons/supervisord.d/bootstrap.conf
[program:bootstrap]
command=/root/bootstrap.sh
redirect_stderr=true
stdout_logfile=/var/log/bootstrap.log
autostart=true
autorestart=false
