# commons/supervisord.d/socks-proxy.conf
[program:socks-proxy]
command=/root/socks-proxy.sh
startsecs=2
stopwaitsecs=10
startretries=30
user=root
redirect_stderr=true
stdout_logfile=/var/log/socks-proxy
autostart=false
